# templates/systemd/waf.service
# WAF systemd service unit
# Install to: /etc/systemd/system/waf.service
# Enable:     systemctl enable --now waf

[Unit]
Description=THE MOTHA FUCKING WAF - Multi-layer Web Application Firewall
Documentation=https://git.omada.cafe/atf/waf
# Start after network is up and Nginx is running
After=network-online.target nginx.service
Wants=network-online.target
PartOf=nginx.service
# If WAF crashes, Nginx will still serve but without protection. PartOf=nginx.service to stops nginx if WAF dies.

[Service]
Type=simple
User=waf
Group=waf

# Binary and config paths
ExecStart=/usr/local/bin/waf -config /etc/waf/config.yaml

# Reload config on SIGHUP (not yet implemented — restart instead)
# ExecReload=/bin/kill -HUP $MAINPID

# Restart policy: always restart on crash, with 5s backoff
Restart=always
RestartSec=5s

# Security Hardening
# Run as a dedicated unprivileged user (no root needed for port > 1024)
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/log/webfirewall
ReadOnlyPaths=/etc/waf

# Capability: none needed (listening on port 7616 > 1024)
CapabilityBoundingSet=
AmbientCapabilities=

# Restrict syscalls to those needed by a Go HTTP server
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Environment
# Load the token secret from a file for more secure than embedding in unit
# Create with: echo -n "$(openssl rand -hex 32)" > /etc/waf/token.secret
# chmod 600 /etc/waf/token.secret
# chown waf:waf /etc/waf/token.secret
EnvironmentFile=-/etc/waf/environment

# Logging
# With output="-" in config.yaml, logs to stdout.
# systemd captures this to journald automatically.
# View logs: journalctl -u waf -f
# View JSON: journalctl -u waf -o json | jq .
StandardOutput=journal
StandardError=journal
SyslogIdentifier=waf

# Resource Limits
# Increase file descriptors for high-connection scenarios
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
